# launchd-plists/ai.hermes.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.hermes.gateway</string>
  <key>UserName</key><string>jakeclaw</string>
  <key>GroupName</key><string>staff</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jakeclaw/.hermes-venv/bin/python</string>
    <string>-m</string>
    <string>hermes_cli.main</string>
    <string>gateway</string>
    <string>run</string>
    <string>--replace</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/jakeclaw/.hermes-venv/lib/python3.13/site-packages</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/Users/jakeclaw/.hermes-venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>VIRTUAL_ENV</key><string>/Users/jakeclaw/.hermes-venv</string>
    <key>HERMES_HOME</key><string>/Users/jakeclaw/.hermes</string>
    <key>HOME</key><string>/Users/jakeclaw</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key><false/>
  </dict>
  <key>ThrottleInterval</key><integer>30</integer>

  <key>StandardOutPath</key><string>/Users/jakeclaw/.hermes/logs/gateway.log</string>
  <key>StandardErrorPath</key><string>/Users/jakeclaw/.hermes/logs/gateway.error.log</string>
</dict>
</plist>
