# launchd-plists/ai.jakeclaw.session_compactor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>ai.jakeclaw.session_compactor</string>
  <key>UserName</key><string>jakeclaw</string>
  <key>GroupName</key><string>staff</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/jakeclaw/workers/bin/session_compactor.sh</string>
  </array>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/jakeclaw/workers/logs/session_compactor.out</string>
  <key>StandardErrorPath</key><string>/Users/jakeclaw/workers/logs/session_compactor.err</string>
  <key>WorkingDirectory</key><string>/Users/jakeclaw</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key><string>/Users/jakeclaw</string>
  </dict>
</dict>
</plist>
